# systemd/medalert.service
[Unit]
Description=medAlert daily medicine check
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=oneshot
ExecStart=/absolute/path/to/medAlert --check